# debian/vplane-controller.socket
[Socket]
ListenStream=/var/run/vyatta/vplaned-config.socket

[Install]
WantedBy=sockets.target
